# curator1-ruksan.conf
[program:curator1-ruksan]
command=python3 /home/curator1/curator1-ruksan/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator1/curator1-ruksan/test.err.log
stdout_logfile=/home/curator1/curator1-ruksan/test.out.log
